# Chatroom/Chatroom/chat-server.conf
[program:chat-server]
command=/server.js
directory=/
autostart=true
autorestart=true
user=your_username
